# jetbot_jupyter.service

[Unit]
Description=Jupyter Notebook Service

[Service]
Type=simple
User=dre
ExecStart=/bin/sh -c "jupyter lab --ip=0.0.0.0 --no-browser"
WorkingDirectory=/home/dre
Restart=always

[Install]
WantedBy=multi-user.target
